# airootfs/etc/systemd/system/install.service
[Unit]
Description=Installation Wizard

[Service]
Type=idle
ExecStart=/root/autorun.sh

[Install]
WantedBy=default.target
